# tests/unit/fake_data_root/storage/ceph-mon/run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
After=blockdev@dev-disk-by\x2dlabel-cloudimg\x2drootfs.target

[Mount]
Where=/
What=/dev/disk/by-label/cloudimg-rootfs
Type=ext4
