# tests/unit/fake_data_root/vault/etc/systemd/system/corosync.service.d/overrides.conf
[Service]
TimeoutStartSec=180
TimeoutStopSec=600