# tests/unit/fake_data_root/openstack/var/lib/juju/agents/unit-neutron-openvswitch-1/charm/templates/kilo/neutron-openvswitch-networking-sriov.conf
description "Configure SRIOV Virtual Functions"
author "Frode Nordahl <frode.nordahl@gmail.com>"

start on virtual-filesystems

task
console log
script
    [ -x "/usr/local/bin/neutron-openvswitch-networking-sriov.sh" ] || exit 0
    exec /usr/local/bin/neutron-openvswitch-networking-sriov.sh start
end script
