# tests/unit/fake_data_root/openstack/run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
After=blockdev@dev-disk-by\x2duuid-51babbe8\x2df78f\x2d46a4\x2d8830\x2dd351c3830325.target

[Mount]
Where=/
What=/dev/disk/by-uuid/51babbe8-f78f-46a4-8830-d351c3830325
Type=ext4
