# tests/unit/fake_data_root/openstack/run/systemd/generator/var-lib-nova-instances.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
After=vaultlocker-decrypt@88244ad9-372d-427e-9d82-c411c73d900a.service
Requires=vaultlocker-decrypt@88244ad9-372d-427e-9d82-c411c73d900a.service
After=blockdev@dev-mapper-crypt\x2d88244ad9\x2d372d\x2d427e\x2d9d82\x2dc411c73d900a.target

[Mount]
Where=/var/lib/nova/instances
What=/dev/mapper/crypt-88244ad9-372d-427e-9d82-c411c73d900a
Type=xfs
Options=defaults,nofail,x-systemd.requires=vaultlocker-decrypt@88244ad9-372d-427e-9d82-c411c73d900a.service,comment=vaultlocker
