# tests/unit/fake_data_root/kubernetes/var/lib/juju/agents/unit-kubernetes-master-0/charm/templates/cdk-service-kicker.service
[Unit]
Description=cdk-service-kicker

[Service]
ExecStart=/usr/bin/cdk-service-kicker
Restart=on-failure
Type=simple

[Install]
WantedBy=multi-user.target
