# tests/unit/fake_data_root/kubernetes/var/lib/juju/agents/unit-kubernetes-master-0/charm/templates/service-iptables-fix.service
[Unit]
Description=Apply iptables rule for KUBE-MARK-DROP
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/kube-proxy-iptables-fix.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
