# tests/unit/fake_data_root/openstack/etc/systemd/system/nova-compute.service.d/99-mount.conf
[Unit]
Requires=var-lib-nova-instances.mount
After=var-lib-nova-instances.mount