# tests/unit/fake_data_root/kubernetes/usr/lib/systemd/user/pk-debconf-helper.service
[Unit]
Description=debconf communication service
Requires=pk-debconf-helper.socket
RefuseManualStart=true

[Service]
ExecStart=/usr/lib/packagekit/pk-debconf-helper
