# tests/unit/fake_data_root/kubernetes/etc/systemd/system/cdk.master.leader.file-watcher.path
[Path]
PathChanged=/root/cdk/basic_auth.csv
PathChanged=/root/cdk/known_tokens.csv
PathChanged=/root/cdk/serviceaccount.key

[Install]
WantedBy=multi-user.target