# tests/unit/fake_data_root/kubernetes/etc/systemd/system/snap.lxd.activate.service
[Unit]
# Auto-generated, DO NOT EDIT
Description=Service for snap application lxd.activate
Requires=snap-lxd-21835.mount
Wants=network.target
After=snap-lxd-21835.mount network.target snapd.apparmor.service
X-Snappy=yes

[Service]
EnvironmentFile=-/etc/environment
ExecStart=/usr/bin/snap run lxd.activate
SyslogIdentifier=lxd.activate
Restart=no
WorkingDirectory=/var/snap/lxd/21835
TimeoutStopSec=30
Type=oneshot
Delegate=true

[Install]
WantedBy=multi-user.target
