# tests/unit/fake_data_root/kubernetes/etc/systemd/system/containerd.service.d/proxy.conf
[Service]
Environment="HTTP_PROXY=" "HTTPS_PROXY=" "NO_PROXY=127.0.0.1,localhost,::1"