# tests/unit/fake_data_root/storage/ceph-mon/etc/systemd/system/snap-snapd-14549.mount
[Unit]
Description=Mount unit for snapd, revision 14549
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/snapd_14549.snap
Where=/snap/snapd/14549
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
