# tests/unit/fake_data_root/openstack/var/lib/juju/agents/unit-neutron-openvswitch-1/charm/templates/os-charm-phy-nic-mtu.conf
description "Apply MTU for Neutron data ports"

start on runlevel [2345]

task

script
    devs="{{ devs }}"
    mtu="{{ mtu }}"
    if [ -n "$mtu" ]; then
        tmpfile=`mktemp`
        echo $devs > $tmpfile
        while read -r dev; do
            [ -n "$dev" ] || continue
            rc=0
            # Try all devices before exiting with error
            ip link set $dev mtu $mtu || rc=$?
        done <  $tmpfile
        rm $tmpfile
        [ $rc = 0 ] || exit $rc
    fi
end script
